# depricated_playbooks/badge_hub.conf
[program:badge_hub]
command=docker-compose up
directory=/home/pi/openbadge-hub-py/
autostart=true
autorestart=unexpected
redirect_stderr=true