# playbooks/badge_hub.conf
[program:badge_hub]
command=/home/pi/openbadge-hub-py/badge_hub.py -m server pull
directory=/home/pi/openbadge-hub-py/
autostart=true
autorestart=unexpected

