# playbooks/report_ip.conf
[program:report_ip]
command=/home/pi/openbadge-hub-py/hub_manager.py
directory=/home/pi/openbadge-hub-py/
user=pi
autostart=true
autorestart=unexpected

